# supervisor.conf

[program:airmash-server]
command = /app/airmash-server
stdout_logfile = syslog
stderr_logfile = syslog

[supervisord]
logfile_maxbytes = 500MB
nodaemon = true
