# supervisor.conf

[program:airmash-server]
command = /app/airmash-server

[supervisord]
logfile_maxbytes = 500MB
nodaemon = true
